# root.d/etc/profile.d/nvim-default-editor.sh
# neovim as default editor
EDITOR=/usr/bin/nvim
